# systemd/ble-agent.service
[Unit]
Description=BlueZ Pairing Agent (Just Works)
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bt-agent -c NoInputNoOutput
Restart=always
RestartSec=1

# Log output to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
